# .local/profile.d/00_bashrc.profile
#!/usr/bin/env bash

# Load bashrc
[ -f "$HOME/.bashrc" ] && . "$HOME/.bashrc"
