# .local/profile.d/30_variables.profile
#!/usr/bin/env bash

# Set gnupg variable
export GNUPGHOME="$XDG_CONFIG_HOME/gnupg"

# Set ICEauthority
export ICEAUTHORITY="$XDG_CONFIG_HOME"/ICEauthority

# Set terminal variable
export TERMINAL=kitty

# Set browser variable
export BROWSER=chromium-browser

# Set editor variable
export EDITOR=nvim
