# .local/profile.d/40_terminal.profile
#!/usr/bin/env bash

export TERMINAL=xfce4-terminal
