# .local/profile.d/docker.profile
#!/usr/bin/env bash

export DOCKER_CONFIG="$XDG_CONFIG_HOME/docker"
export DOCKER_CLI_EXPERIMENTAL=enabled
