# .local/profile.d/rust.profile
#!/usr/bin/env bash

export RUSTUP_HOME="$HOME/.local/share/rustup"
export CARGO_HOME="$HOME/.local/share/cargo"
source "$HOME/.local/share/cargo/env"
