# .local/profile.d/screen_power_saving.sh
# turn off screensaver
xset s off

# turn off DPMS
xset -dpms
