# .local/profile.d/scripts.sh
export PATH="$(find -L $HOME/.scripts -type d | tr '\n' ':')$PATH"
